# templates/photo_mover.service.j2
[Unit]
Description=scheduled photo mover

[Service]
User={{ photo_mover_user }}
Group={{ photo_mover_group }}
ExecStart=/opt/photo_mover/photo_mover.sh
